# crontab
*/5 * * * * /usr/bin/php /www/SSPanel-Uim/xcat  Cron
