# crontab
*/5 * * * * php /www/SSPanel-Uim/xcat  Cron
